# deployment/config/typemaster.service
[Unit]
Description=TypeMaster Typing Tutor — Spring Boot Application
Documentation=https://github.com/your-org/typing-tutor
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=typemaster
Group=typemaster

# Application location
WorkingDirectory=/opt/typemaster/app

# Load environment variables from file
EnvironmentFile=/etc/typemaster/typemaster.env

# JVM tuning for 12 GB RAM VM (adjust Xmx if needed)
ExecStart=/usr/bin/java \
    -Dspring.profiles.active=prod \
    -Xms256m \
    -Xmx512m \
    -XX:+UseG1GC \
    -XX:MaxGCPauseMillis=200 \
    -Djava.security.egd=file:/dev/./urandom \
    -jar /opt/typemaster/app/app.jar

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/typemaster/app /var/log/typemaster

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=typemaster

[Install]
WantedBy=multi-user.target
